# dist/usbguard.service
[Unit]
Description=USBGuard daemon
Wants=systemd-udevd.service local-fs.target

[Service]
Type=simple
ExecStart=/usr/sbin/usbguard-daemon -k -c /etc/usbguard/usbguard-daemon.conf
Restart=on-failure

[Install]
WantedBy=basic.target
