# src/DBus/usbguard-dbus.service.in
[Unit]
Description=USBGuard D-Bus Service
Requires=usbguard.service
Documentation=man:usbguard-dbus(8)

[Service]
Type=dbus
BusName=org.usbguard1
ExecStart=%{sbindir}%/usbguard-dbus --system
Restart=on-failure

[Install]
WantedBy=multi-user.target
Alias=dbus-org.usbguard.service

